# v874-supervisor.conf
;
; added cst 20160930
;
; auto-restart redddate gunicorn/daphne startscript.
;
; Linked from: /etc/supervisor/conf.d/dtr4.conf
;
;   ~# supervisorctl stop dtr4
;   ~# supervisorctl start dtr4
;
[program:dtr4]
command = /var/www/www.elligue.com/dtr4/v874-gunicorn-start.sh
; command = /var/www/www.elligue.com/dtr4/v874-daphne-start.sh
; stdout_logfile = /var/log/supervisor/dtr4_daphne.log
stdout_logfile = /var/log/supervisor/dtr4_gunicorn.log
user = cst
redirect_stderr = true
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8
