# v874-supervisor.conf
; last modified: cst 20161012
;
; auto-restart El Ligue (dtr4) gunicorn/daphne/workers control scripts.
;
; Linked from: /etc/supervisor/conf.d/dtr4.conf
;
;
; To control the *SGI (Gunicorn, Daphne, etc) use:
;
;   ~# supervisorctl stop dtr4-gateway
;   ~# supervisorctl start dtr4-gateway
;
; To control Django worker when running a ASGI gateway, use:
;
;   ~# supervisorctl stop dtr4-workers
;   ~# supervisorctl start dtr4-workers
;
[program:dtr4-gateway]
command = /var/www/www.elligue.com/dtr4/v874-start-gateway.sh
stdout_logfile = /var/log/supervisor/dtr4-gateway.log
user = cst
redirect_stderr = true
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8
;
;
;
[program:dtr4-workers]
command = /var/www/www.elligue.com/dtr4/v874-start-workers.sh
stdout_logfile = /var/log/supervisor/dtr4-workers.log
user = cst
redirect_stderr = true
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8
